# packages/backend/systemd/buffpenguin.service
[Unit]
Description=BuffPenguin Gym Tracker Backend
After=network.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/BuffPenguin/packages/backend
ExecStart=/usr/bin/node dist/index.js
Restart=on-failure
RestartSec=5
Environment=NODE_ENV=production
Environment=PORT=3000
Environment=DB_PATH=/home/pi/BuffPenguin/packages/backend/data/buffpenguin.db

[Install]
WantedBy=multi-user.target
